# bin/my_installer/containers_installer/Extra/mullvad/app/conf.d/amixer.conf
[program:amixer]
command=amixer set Master playback 100%
autorestart=true
