# bin/my_installer/containers_installer/Extra/mullvad/app/conf.d/mullvad-browser.conf
[program:mullvad-browser]
command=mullvad-browser
autorestart=true
