# bin/my_installer/containers_installer/Extra/mullvad/app/conf.d/pulseaudio.conf
[program:pulseaudio]
command=pulseaudio
autorestart=true
