# bin/my_installer/containers_installer/Extra/mullvad/app/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 4433 localhost:5900
autorestart=true
